# root/etc/uci-defaults/80_domain-proxy
#!/bin/sh

touch /etc/config/domain-proxy
uci add domain-proxy proxy
uci commit

return 0